# mailman3-web-weekly.service
[Unit]
Description=GNU Mailing List Manager 3 Web UI / runjobs weekly

# run only if main service is running
BindsTo=mailman3-web.service

[Service]
Type=oneshot
ExecStart=@BINDIR@/mailman-web runjobs weekly
User=@MMUSER@
Group=@MMGROUP@
