# mailman3-web-hourly.service
[Unit]
Description=GNU Mailing List Manager 3 Web UI / runjobs hourly

# run only if main service is running
BindsTo=mailman3-web.service

[Service]
Type=oneshot
ExecStart=@BINDIR@/mailman-web runjobs hourly
User=@MMUSER@
Group=@MMGROUP@
